# config/prqld.service
[Unit]
Description=PrQL Server

[Service]
Type=simple
ExecStart=/usr/bin/prqld
Restart=always

[Install]
WantedBy=multi-user.target
